# etc/supervisord/supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=apache2ctl -DFOREGROUND
stopsignal=WINCH
stopwaitsecs=28
user=root
killasgroup=true
stopasgroup=true

